# deploy/com.pipboi.watchdog.plist
<?xml version="1.0" encoding="UTF-8"?>
<!-- Vault-bot fleet watchdog — runs fleet_watchdog.py every 6h.
     LaunchAgent (not cron): runs in your user session and fires a missed job on wake,
     which matters on a laptop that sleeps. Reclaims empty-ATA rent, auto-sells orphans
     worth >= 0.05 SOL, and alerts on frozen bots / stalled sidecar / dead loops / a
     fresh ghost-close streak. fleet_watchdog.py holds a single-instance lock, so this
     agent and a running watchdog_loop.sh cannot double-run / double-sell.

     Install:
       Install with: ./deploy/install.sh  (fills the real path + loads the agent)
       launchctl bootstrap gui/$(id -u) ~/Library/LaunchAgents/com.pipboi.watchdog.plist
       launchctl enable gui/$(id -u)/com.pipboi.watchdog
     Remove:
       launchctl bootout gui/$(id -u)/com.pipboi.watchdog
     Run once now (test):
       launchctl kickstart -k gui/$(id -u)/com.pipboi.watchdog
-->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.pipboi.watchdog</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>-lc</string>
    <string>cd __PIPBOI_DIR__/src &amp;&amp; python3 fleet_watchdog.py</string>
  </array>
  <key>WorkingDirectory</key>
  <string>__PIPBOI_DIR__/src</string>
  <key>RunAtLoad</key>
  <true/>
  <key>StartInterval</key>
  <integer>21600</integer>
  <key>StandardOutPath</key>
  <string>__PIPBOI_DIR__/src/logs/watchdog_launchd.log</string>
  <key>StandardErrorPath</key>
  <string>__PIPBOI_DIR__/src/logs/watchdog_launchd.log</string>
</dict>
</plist>
